# root.alpine/etc/services.d/smbd/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/smbd --foreground --no-process-group
