# root.alpine/etc/services.d/dbus/run
#!/usr/bin/execlineb -P

exec /usr/bin/dbus-daemon --system --nofork
