# host/systemd/podman-compose@libvirt-pod.service.d/override.conf
[Service]
Environment=DIST=alpine
#Environment=DIST=fedora
