# root.fedora/usr/lib/systemd/system/libvirtd-ro.socket.d/10-container-overrides.conf
[Unit]
After=container-init.service

[Socket]
SocketMode=0660
SocketGroup=libvirt_access
