# root.fedora/usr/lib/systemd/system/smb.service.d/10-container-overrides.conf
[Unit]
After=container-init.service
ConditionEnvironment=SMBD=1
